# debian/postinst
#!/bin/sh

set -e

update-rc.d package defaults
