# debian/postinst
#!/bin/sh

set -e

update-rc.d iptables-fsinf defaults

#DEBHELPER#
